# pi-control/telegram/claude-telegram-bot.service
[Unit]
Description=Claude Code Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aarzner
WorkingDirectory=/home/aarzner/home-infra
ExecStart=/home/aarzner/home-infra/pi-control/telegram/.venv/bin/python /home/aarzner/home-infra/pi-control/telegram/bot.py
Restart=on-failure
RestartSec=10
Environment=PATH=/home/aarzner/.local/bin:/home/aarzner/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
