# pi-flag-cam/pi/config/crash-monitor.service
[Unit]
Description=Pi Flag Cam crash monitor (logs to boot partition)
After=local-fs.target

[Service]
Type=simple
User=root
ExecStart=/home/aarzner/pi-flag-cam/crash-monitor.sh
Restart=always
RestartSec=5
# Hardening (root is required to remount /boot/firmware rw for the crash log).
ProtectHome=read-only
ReadWritePaths=/boot/firmware /run
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
