# pi-flag-cam/pi/config/pi-flag-cam.service
[Unit]
Description=Pi Flag Cam API Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=aarzner
ExecStart=/usr/bin/python3 /home/aarzner/pi-flag-cam/server.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

# Hardening
ProtectHome=read-only
ProtectSystem=strict
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/tmp
# HID access to the Luxafor without world-writable device nodes.
SupplementaryGroups=plugdev
# Firewall (systemd cgroup BPF — works even though iptables is blacklisted):
# only the operator laptop (t14, pinned by DHCP reservation) + localhost may
# reach :8080. Keep the IP in sync with PI_FLAG_CAM_ALLOW_IP in config.sh.
IPAddressAllow=localhost
IPAddressAllow=192.168.86.86
IPAddressDeny=any

[Install]
WantedBy=multi-user.target
