# pi-control/monitoring/config/victoriametrics.service
[Unit]
Description=VictoriaMetrics single-node
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/victoria-metrics \
    -storageDataPath=/var/lib/victoria-metrics \
    -retentionPeriod=6 \
    -promscrape.config=/etc/victoriametrics/scrape.yml \
    -httpListenAddr=:8428
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
